# templates/systemd/mrs.service.j2
{#
SPDX-FileCopyrightText: 2023 Nikita Chernyi

SPDX-License-Identifier: AGPL-3.0-or-later
#}

[Unit]
Description=Matrix Rooms Search API ({{ mrs_identifier }})
{% for service in mrs_systemd_required_systemd_services_list %}
Requires={{ service }}
After={{ service }}
{% endfor %}
DefaultDependencies=no

[Service]
Type=simple
Environment="HOME={{ devture_systemd_docker_base_systemd_unit_home_path }}"
ExecStartPre=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} kill {{ mrs_identifier }} 2>/dev/null'
ExecStartPre=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} rm {{ mrs_identifier }} 2>/dev/null'

ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} create \
			--rm \
			--name={{ mrs_identifier }} \
			--log-driver=none \
			--network={{ mrs_container_network }} \
			--user={{ mrs_uid }}:{{ mrs_gid }} \
			--cap-drop=ALL \
			--read-only \
			--hostname={{ mrs_hostname }} \
			{% if mrs_container_http_bind_port %}
			-p {{ mrs_container_http_bind_port }}:8080 \
			{% endif %}
			--ulimit nofile=262144:262144 \
			--label-file={{ mrs_base_path }}/labels \
			--mount type=bind,src={{ mrs_base_path }},dst=/mnt \
			--tmpfs=/tmp:rw,noexec,nosuid,size=512m \
			{{ mrs_container_image }} \
			-c {{ mrs_path_config }}

{% for network in mrs_container_additional_networks %}
ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} network connect {{ network }} {{ mrs_identifier }}
{% endfor %}

ExecStart={{ devture_systemd_docker_base_host_command_docker }} start --attach {{ mrs_identifier }}

ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} kill {{ mrs_identifier }} 2>/dev/null'
ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} rm {{ mrs_identifier }} 2>/dev/null'
Restart=always
RestartSec=30
SyslogIdentifier={{ mrs_identifier }}

[Install]
WantedBy=multi-user.target
